# vote-secret-systemd-fixed.service
[Unit]
Description=Vote Secret v2.0 Backend Service
After=network.target mongodb.service
Wants=mongodb.service

[Service]
Type=exec
User=vote-secret
Group=vote-secret
WorkingDirectory=/opt/vote-secret/backend
Environment=PATH=/opt/vote-secret/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/opt/vote-secret/backend
ExecStart=/opt/vote-secret/venv/bin/gunicorn --config /opt/vote-secret/config/gunicorn.conf.py server:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vote-secret /var/log/vote-secret
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vote-secret

[Install]
WantedBy=multi-user.target